# crontab
* * * * * root sh /app/generate-docs-manifest.sh >> /var/log/cron.log
