# crontab
* * * * * root sh /app/process-docs.sh >> /var/log/cron.log
